# meshtastic-mqtt-relay.service
# Example systemd service file for meshtastic-mqtt-relay
# Copy this file to /etc/systemd/system/meshtastic-mqtt-relay.service
# Edit the configuration parameters below to match your setup
# Enable and start with:
#   sudo systemctl enable meshtastic-mqtt-relay
#   sudo systemctl start meshtastic-mqtt-relay

[Unit]
Description=Meshtastic MQTT Binary-to-JSON Relay
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=nobody
Group=nogroup
ExecStart=/usr/local/bin/meshtastic-mqtt-relay \
  -source-broker tcp://mqtt.meshtastic.org:1883 \
  -source-topic "msh/US/2/json/#" \
  -dest-broker tcp://localhost:1883 \
  -dest-topic meshtastic/json

# Restart on failure
Restart=on-failure
RestartSec=10s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meshtastic-mqtt-relay

[Install]
WantedBy=multi-user.target
